# prod/midas_protocol/deployment/systemd/magellan-midas-protocol.service
[Unit]
Description=Magellan MIDAS Protocol Strategy (Account PA3DDLQCBJSE)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ssm-user
Group=ssm-user
WorkingDirectory=/home/ssm-user/magellan

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="STRATEGY=midas_protocol"
Environment="ALPACA_ACCOUNT_ID=PA3DDLQCBJSE"
Environment="ENVIRONMENT=production"

# Execution
ExecStart=/home/ssm-user/magellan/.venv/bin/python /home/ssm-user/magellan/prod/midas_protocol/runner.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=magellan-midas-protocol

# Resource limits
MemoryLimit=2G
CPUQuota=100%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
